# xxx.conf
[program:xxx]
command=python3 /home/curator2/xxx/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/xxx/test.err.log
stdout_logfile=/home/curator2/xxx/test.out.log